# profile.d/path.sh
#!/usr/bin/bash.exe
PATH=/c/Users/Public/OneDrive/Documents/Projects/BYU-I/trrmann.github.io/scripts:${PATH}