# profile.d/trrmann.github.io.sh
#!/usr/bin/bash.exe
echo "Load trrmann.github.io profile!"